# deploy_ng/roles/base/templates/get_info.service.j2
[Unit]
Description=get_info
Requires=network.target
After=syslog.target network.target
[Service]
ExecStart=/usr/bin/get_info.sh -s
User=andreas1
Group=andreas1


[Install]
WantedBy=multi-user.target